# centos/6.6/supervisord/supervisord-ds389.conf
[program:ds389]
command=/bin/bash -c "exec /usr/sbin/ns-slapd -D /etc/dirsrv/slapd-${SERVER_IDENTIFIER} -d stats+stats2+backend -i /var/run/dirsrv/slapd-${SERVER_IDENTIFIER}.pid -w /var/run/dirsrv/slapd-${SERVER_IDENTIFIER}.startpid"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:ds389-admin]
command=/bin/bash -c "exec /usr/sbin/httpd.worker -k start -DFOREGROUND -f /etc/dirsrv/admin-serv/httpd.conf"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
